# xsession/lock.service
[Unit]
Description=Lock the screen on resume from suspend

[Service]
User=pkpao
Environment=DISPLAY=:0
ExecStart=/usr/local/bin/slock

[Install]
WantedBy=suspend.target
